# deploy/kuber.service
[Unit]
Description=Kuber API (Express + SQLite)
After=network.target

[Service]
Type=simple
User=kuber
WorkingDirectory=/opt/kuber/server
EnvironmentFile=/opt/kuber/server/.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=3
# Мягкий лимит памяти-подсказка (better-sqlite3 синхронный, много не ест).
Environment=NODE_OPTIONS=--max-old-space-size=512

[Install]
WantedBy=multi-user.target
